# roles/hardening/files/suricata.target
[Unit]
Description=Suricata Target

[Install]
WantedBy=multi-user.target
