# roles/hardening/files/suricata-update.service
[Unit]
Description=Suricata update rules

[Service]
ExecStart=suricata-update
Restart=always

[Install]
WantedBy=suricata.target
